# services/sscrt-sxmr-bot.service
[Unit]
Description=sscrt-sxmr pair arbitrage bot

[Service]
Type=simple
ExecStart=python3 arb_bot_2p.py sscrt-sxmr-config
Restart=on-failure
WorkingDirectory=/root/arb/src

[Install]
WantedBy=multi-user.target